# systemd/horus-cli.service
[Unit]
Description=BeoutOS CLI Service - Secure Management Console
After=horus-overlay.service
After=horus-activation.service
After=horus-lockdown.service
ConditionPathExists=/mnt/horus-config/license/activated

[Service]
Type=simple
ExecStart=/usr/sbin/horus-cli
Restart=on-failure
RestartSec=3
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

Environment=HORUS_MODE=active
Environment=BEOUTOS_CONFIG_MOUNT=/mnt/horus-config
Environment=HORUS_LICENSE=activated

ExecStartPre=/bin/sh -c 'systemctl mask getty@tty1.service 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
Alias=getty@tty1.service
